# systemd/dumang-sync.service
[Unit]
Description=DuMang DK6 Layer Sync

[Service]
Type=Simple
ExecStart=dumang-sync

[Install]
WantedBy=multi-user.target
