# systemd/dumang-sync-python.service
[Unit]
Description=DuMang DK6 Layer Sync (Python Invoke)

[Service]
Type=Simple
ExecStart=dumang-sync

[Install]
WantedBy=multi-user.target
